# praktikum-modul-1-f01/task-2/crontabs
0 */1 * * * cd /home/ardhikakrsh/Documents/Kuliah/Sisop/Praktikum-1/praktikum-modul-1/task-2 && ./yu_database.sh

